# src/main/packaging/scripts/control/deb/postrm
#!/bin/sh

update-rc.d -f debhome remove